# overlay/etc/systemd/system/netctl@bridge.service
.include /usr/lib/systemd/system/netctl@.service

[Unit]
Description=Docker bridge
After="netctl@main.service"
